# supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
logfile_maxbytes=5MB

[program:bot]
command=python -u main.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/bot_stdout.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/app/logs/bot_stderr.log
stderr_logfile_maxbytes=5MB

[program:dashboard]
command=python -u run_dashboard.py --no-browser
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/dashboard_stdout.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/app/logs/dashboard_stderr.log
stderr_logfile_maxbytes=5MB
